# cfg/systemd/system/pg_vacuum.timer
[Unit]
Description=Run pg_vacuum daily

[Timer]
OnCalendar=*-*-* 23:40:00
Persistent=true

[Install]
WantedBy=timers.target
