# cfg/systemd/system/pgdump_cleanup.timer
[Unit]
Description=Run pgdump_cleanup daily

[Timer]
OnCalendar=*-*-* 23:35:00
Persistent=true

[Install]
WantedBy=timers.target
