# docs/supervisord/tcp_mikeandwan.us.ini
[program:mikeandwan.us]
command=/usr/local/bin/dotnet /srv/www/mikeandwan.us/mikeandwan.us.dll
directory=/srv/www/mikeandwan.us/
autostart=true
autorestart=true
stderr_logfile=/var/log/mikeandwan.us.err.log
stdout_logfile=/var/log/mikeandwan.us.out.log
environment=ASPNETCORE_ENVIRONMENT="production",MAW_Environment__DbConnectionString="Server=localhost;Port=5432;Database=maw_website;User Id=svc_www_maw;Password=???",MAW_Email__User="email@example.com",MAW_Email__Password="???",MAW_Google__RecaptchaSiteKey="???",MAW_Google__RecaptchaSecretKey="???"
user=svc_www_maw
stopsignal=INT
stopasgroup=true
killasgroup=true
